# deploy-ecs/config-docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:laravel-job]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --timeout=3600 
autostart=true
autorestart=true    
numprocs=10
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stderr
stopwaitsecs=3600

# [program:laravel-retry]
# process_name=%(program_name)s_%(process_num)02d
# command=php /var/www/html/artisan queue:retry --range=295-305
# autostart=true
# autorestart=false
# numprocs=1
# stderr_logfile=/dev/stderr
# stdout_logfile=/dev/stderr
# stopwaitsecs=3600

# [program:laravel-schedule]
# process_name=%(program_name)s_%(process_num)02d
# command=php /var/www/html/artisan schedule:run
# autostart=true
# autorestart=false
# numprocs=1
# stderr_logfile=/dev/stderr
# stdout_logfile=/dev/stderr
# stopwaitsecs=3600
